# docker/musicbrainz-tests/service/website/run
#!/usr/bin/env bash

exec unbuffer "$MBS_ROOT"/docker/scripts/start_musicbrainz_server.sh
