# admin/nginx/service/log/run
#!/bin/bash
mkdir -p ./main
exec multilog t \
       +'*' s1000000 n50 ./main \
