# docker/musicbrainz-tests/service/mb-mail-service/run
#!/usr/bin/env bash

export APP_SMTP_PORT=1025

exec mb-mail-service
