# docker/musicbrainz-tests/service/mb-mail-service/run
#!/usr/bin/env bash

exec mb-mail-service
