# docker/musicbrainz-tests/service/vnu/run
#!/bin/sh

exec java -cp /home/musicbrainz/vnu/vnu.jar \
    nu.validator.servlet.Main 8888 2>&1
