# docker/musicbrainz-tests/service/mailpit/run
#!/bin/sh

exec chpst -u musicbrainz:musicbrainz mailpit
