# docker/musicbrainz-json-dump/crontab
SHELL=/bin/bash
30 * * * * $HOME/musicbrainz-server/admin/run admin/cron/hourly-json-dump.sh
10 0 * * * $HOME/musicbrainz-server/admin/run admin/cron/daily-json-dump.sh
