# docker/musicbrainz-tests/service/smtp/run
#!/bin/sh

exec python3.13 -m aiosmtpd -n -c DebuggingServer localhost:25
