# docker/musicbrainz-search-indexes-dump/crontab
SHELL=/bin/bash
10 4 * * * $HOME/musicbrainz-server/admin/run admin/cron/daily-search-indexes-dump.sh
