# docker/musicbrainz-solr-backup/crontab
SHELL=/bin/bash
10 4 * * * $HOME/musicbrainz-server/admin/run admin/cron/daily-solr-backup.sh
